# services/supervisord/supervisord.d/confd.ini
[program:confd]
user=root
command=confd -backend nacos -node http://127.0.0.1:8848 -watch -onetime=false -log-level="debug" -namespace="public"

autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/confd/confd_out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10